# src/server/init/kibana4.conf
description "Kibana4"
author "Philip Hutchins"

start on (filesystem and net-device-up)
stop on runlevel [!2345]

respawn
respawn limit 5 30
limit nofile 65550 65550

chdir /opt/kibana/current

setuid kibana
setgid kibana

script
  export PATH="$PATH:$HOME/.rvm/bin"
  export KIBANA_HOME="/opt/kibana/current"

  exec $KIBANA_HOME/bin/kibana
end script

emits kibana-server-running
