# src/dev/perf_test/benchmark/logstash-7.5.1/vendor/bundle/jruby/2.5.0/gems/pleaserun-0.0.30/templates/systemd-user/default/program.service
[Unit]
{{#description}}
Description={{{ description }}}
{{/description}}
 
[Service]
Type=simple
EnvironmentFile={{{home}}}/.config/systemd/user/{{{name}}}.environment
ExecStart={{{program}}} {{{shell_args}}}
Restart=always
WorkingDirectory={{{chdir}}}
StandardOutput=journal

[Install]
WantedBy=multi-user.target
